# etc/systemd/system/dhcpcd.service.d/override.conf
[Unit]
Requires=wpa_supplicant.service
After=wpa_supplicant.service
