# etc/systemd/system/dnsmasq.service.d/override.conf
[Unit]
Requires=wpa_supplicant.service dhcpcd.service hostapd.service systemd-networkd.service
After=wpa_supplicant.service dhcpcd.service hostapd.service systemd-networkd.service
