# etc/systemd/system/hostapd.service.d/override.conf
[Unit]
Requires=sys-subsystem-net-devices-wifi_int.device wpa_supplicant.service dhcpcd.service
After=sys-subsystem-net-devices-wifi_int.device wpa_supplicant.service dhcpcd.service
