# etc/systemd/system/wpa_supplicant.service.d/override.conf
[Unit]
After=sys-subsystem-net-devices-wifi_ext.device
Requires=sys-subsystem-net-devices-wifi_ext.device
